# .ddev/web-build/sshd-supervisor.conf
#ddev-generated
[program:sshd]
command=/usr/local/bin/ai-ssh-setup.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.log
